# backend/systemd/hivemind.timer
[Unit]
Description=Run hivemind every 5 min

[Timer]
OnActiveSec=5min

[Install]
WantedBy=timers.target
